# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
priority=10

[program:backend]
command=/bin/bash -c "cd /app/backend && /app/backend/entrypoint.sh uvicorn app.main:app --host 127.0.0.1 --port 8000"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_err.log
priority=20
environment=PYTHONUNBUFFERED="1"
